# ac-offline.service
[Unit]
Description=AC offline

[Service]
Type=oneshot
RemainAfterExit=yes

# set parameters for battery mode
ExecStart=/usr/local/bin/ac-offline true

# set parameters for AC mode
ExecStop=/usr/local/bin/ac-offline false
